# apps/mybookkeeper/deploy/dramatiq-scheduler.service
[Unit]
Description=MyBookkeeper Gmail Sync Scheduler
After=network.target postgresql.service

[Service]
User=deploy
WorkingDirectory=/srv/mybookkeeper/backend
EnvironmentFile=/srv/mybookkeeper/backend/.env
ExecStart=/srv/mybookkeeper/backend/.venv/bin/python -m app.workers.scheduler_worker
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
