# apps/mybookkeeper/deploy/uvicorn.service
[Unit]
Description=MyBookkeeper FastAPI
After=network.target postgresql.service

[Service]
User=deploy
WorkingDirectory=/srv/mybookkeeper/backend
EnvironmentFile=/srv/mybookkeeper/backend/.env
ExecStart=/srv/mybookkeeper/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2 --root-path /api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
