# apps/myjobhunter/deploy/myjobhunter-backup.timer
[Unit]
Description=MyJobHunter daily database backup at 2 AM

[Timer]
# Run daily at 02:00 local time. Persistent=true so a backup runs on
# the next boot if the VPS was offline at 02:00 (e.g., post-reboot,
# post-deploy).
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
